# systemd/nasa-tunnel.service
[Unit]
Description=NASA reverse SSH tunnel to VPS (bypasses CGNAT)
Documentation=https://github.com/AlexeyBorovskoy/Nasa_home
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
EnvironmentFile=/opt/nasa/config/.env
ExecStart=/usr/bin/autossh -N \
    -R 18080:localhost:8080 \
    -R 12283:localhost:2283 \
    -R 18090:localhost:8090 \
    -R 10022:localhost:22 \
    -o ServerAliveInterval=30 \
    -o ServerAliveCountMax=3 \
    -o ExitOnForwardFailure=yes \
    -o StrictHostKeyChecking=accept-new \
    -i "${VPS_SSH_KEY}" \
    "${VPS_USER}@${VPS_HOST}"
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target