# systemd/nasa-usb-monitor.service
[Unit]
Description=NASA USB Error Monitor (RTL9210B-CG real-time dmesg watcher)
Documentation=https://github.com/AlexeyBorovskoy/Nasa_home
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=10s
ExecStart=/usr/local/sbin/nasa-usb-monitor.sh
EnvironmentFile=-/etc/nasa-monitor/watchdog.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nasa-usb-monitor

[Install]
WantedBy=multi-user.target
