# systemd/nasa-backup.timer
[Unit]
Description=Run NASA Home Cloud database backup daily at 03:00
Documentation=https://github.com/AlexeyBorovskoy/Nasa_home

[Timer]
# 03:00 — before 09:00 Telegram report, off-peak for Jetson
OnCalendar=*-*-* 03:00:00
# Run on next boot if the last trigger was missed (Jetson was off)
Persistent=true
# Spread load ±15 min
RandomizedDelaySec=15m
Unit=nasa-backup.service

[Install]
WantedBy=timers.target
