# systemd/nasa-backup.service
[Unit]
Description=NASA Home Cloud — database backup (pg_dump)
Documentation=https://github.com/AlexeyBorovskoy/Nasa_home
After=local-fs.target network-online.target docker.service
Wants=network-online.target
Requires=docker.service
RequiresMountsFor=/mnt/storage
ConditionPathIsMountPoint=/mnt/storage

[Service]
Type=oneshot
# Adjust NASA_PROJECT_DIR if repo is not at /home/admin/nasa
Environment=NASA_PROJECT_DIR=/home/admin/nasa
ExecStartPre=/bin/mountpoint -q /mnt/storage
ExecStart=/bin/bash "${NASA_PROJECT_DIR}/scripts/backup/backup_databases.sh"
User=admin
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nasa-backup
