# systemd/nasa-sd-wear.service
[Unit]
Description=NASA microSD Wear Check
After=network.target

[Service]
Type=oneshot
ExecStart=/home/admin/nasa/scripts/diagnostics/sd_wear_check.sh
