# systemd/nasa-ssd-recovery.service
[Unit]
Description=NASA SSD Hotplug Auto-Recovery
Documentation=https://github.com/AlexeyBorovskoy/Nasa_home
# Run outside normal boot dependency tree — triggered by udev at any time
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/home/admin/nasa/scripts/storage/ssd_hotplug_recovery.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nasa-ssd-recovery
# Allow up to 120s for mount + preflight + Docker start
TimeoutStartSec=120
