# systemd/nasa-jms583-health.timer
[Unit]
Description=NASA JMS583 USB SSD Health Check — hourly
Documentation=https://github.com/AlexeyBorovskoy/Nasa_home

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
