# systemd/nasa-usb-watchdog.service
[Unit]
Description=NASA USB SSD Recovery Watchdog
After=network.target
# Не запускать одновременно с монтированием — дать ему шанс
After=jetson-nas-mount.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nasa-usb-watchdog.sh
# Читаем Telegram credentials из env файла если есть
EnvironmentFile=-/etc/nasa-monitor/watchdog.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nasa-usb-watchdog
# Разрешаем reboot
User=root
