# systemd/nasa-jms583-health.service
[Unit]
Description=NASA JMS583 USB SSD Health Check
Documentation=https://github.com/AlexeyBorovskoy/Nasa_home
After=network.target nasa-ssd-recovery.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/jms583_health.sh
EnvironmentFile=-/etc/nasa-monitor/telegram.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nasa-jms583

[Install]
WantedBy=multi-user.target
