# systemd/docker.service.d/10-nasa-storage.conf
[Unit]
Wants=jetson-nas-mount.service
After=jetson-nas-mount.service local-fs.target
RequiresMountsFor=/mnt/storage
